# deployment/systemd/wmata-collector.service
[Unit]
Description=WMATA Dashboard Data Collector
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=wmata
Group=wmata
WorkingDirectory=/home/wmata/wmata-dashboard
Environment="PATH=/home/wmata/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/wmata/wmata-dashboard/.env

# Run the continuous combined collector script.
# Invoke the venv interpreter directly (not `uv run`) so the unit doesn't need
# to write a uv cache under the read-only home, and to drop the ~36 MB uv
# wrapper process.
ExecStart=/home/wmata/wmata-dashboard/.venv/bin/python3 scripts/continuous_combined_collector.py

# Restart policy
Restart=always
RestartSec=5s

# Logging
StandardOutput=append:/var/log/wmata/collector.log
StandardError=append:/var/log/wmata/collector-error.log
SyslogIdentifier=wmata-collector

# Resource limits (2 GB RAM instance + 4 GB swap). Measured baseline ~210–290 MB
# RSS; raised from 600M (too tight — transient peaks hit the ceiling, causing
# dropped heartbeat ticks) to 1G to give real headroom above the ~290M peak.
MemoryMax=1G
CPUQuota=50%

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
# The collector writes a PID file under logs/ and the JSONL cold-archive under
# archive/raw_snapshots/ — both must be writable despite ProtectHome=read-only.
ReadWritePaths=/var/log/wmata /home/wmata/wmata-dashboard/logs /home/wmata/wmata-dashboard/archive

[Install]
WantedBy=multi-user.target
